# daemon/godelion-daemon.service
[Unit]
Description=Godelion Daemon Service
After=network.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/godelion
ExecStart=/root/godelion/bin/godelion-daemon
Restart=always
RestartSec=5
StandardOutput=append:/root/godelion/logs/daemon-stdout.log
StandardError=append:/root/godelion/logs/daemon-stderr.log

[Install]
WantedBy=multi-user.target
